# debian/postinst
#!/bin/bash
set -e

# Post-installation script for netpulse

case "$1" in
    configure)
        # Create necessary directories
        mkdir -p /var/lib/netpulse /var/log/netpulse /etc/netpulse
        
        # Set proper ownership for netpulse user
        chown -R netpulse:netpulse /var/lib/netpulse /var/log/netpulse
        chown root:root /etc/netpulse
        
        # Set proper permissions
        chmod 755 /var/lib/netpulse /var/log/netpulse
        chmod 755 /etc/netpulse
        
        # Ensure database directory is writable by netpulse user
        chmod 775 /var/lib/netpulse
        find /var/lib/netpulse -name "*.db" -exec chmod 664 {} \; 2>/dev/null || true
        
        # Ensure log directory is writable by netpulse user
        chmod 775 /var/log/netpulse
        
        # Install configuration file if it doesn't exist
        if [ ! -f /etc/netpulse/netpulse.conf ]; then
            cp /usr/share/doc/netpulse/default.conf /etc/netpulse/netpulse.conf 2>/dev/null || \
            cp /usr/share/netpulse/default.conf /etc/netpulse/netpulse.conf 2>/dev/null || \
            echo "No default config found, creating empty one" && \
            touch /etc/netpulse/netpulse.conf
        fi
        
        # Enable and start services
        systemctl daemon-reload
        systemctl enable netpulse.timer
        systemctl start netpulse.timer
        
        # Enable and start web service
        systemctl enable netpulse-web.service
        systemctl restart netpulse-web.service
        ;;
        
    abort-upgrade|failed-upgrade)
        # Stop services on upgrade failure
        systemctl stop netpulse.timer || true
        systemctl stop netpulse.service || true
        systemctl stop netpulse-web.service || true
        ;;
esac

exit 0
